# ansible/templates/download-pull.timer.j2
[Unit]
Description=Run download pull sync periodically

[Timer]
OnBootSec=2min
OnUnitActiveSec={{ sync_interval }}

[Install]
WantedBy=timers.target
