# ansible/templates/snapraid-sync.service.j2
[Unit]
Description=Snapraid sync
After=local-fs.target

[Service]
Type=oneshot
ExecStart=/usr/bin/snapraid sync
