# ansible/templates/snapraid-sync.timer.j2
[Unit]
Description=Run snapraid sync weekly

[Timer]
OnCalendar=weekly
Persistent=true

[Install]
WantedBy=timers.target
