# ansible/templates/torrent-push.timer.j2
[Unit]
Description=Run torrent push sync periodically

[Timer]
OnBootSec=2min
OnUnitActiveSec={{ sync_interval }}

[Install]
WantedBy=timers.target
